# pure-scripts/chroot_scripts/init_file
#!/bin/bash
### BEGIN INIT INFO
# Provides:          ltsp-ssh-init
# Required-Start:    $network
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:
# Short-Description: Gera chaves SSH e reinicia o serviço
### END INIT INFO

log() {
  echo "[ltsp-ssh-init] $1"
}

case "$1" in
  start)
    log "Gerando chaves SSH (se necessário)..."
    ssh-keygen -A

    log "Reiniciando o serviço SSH..."
    systemctl restart ssh || service ssh restart

    if command -v vmware-modconfig &>/dev/null; then
      log "Recompilando kernel modules do VMware..."
      vmware-modconfig --console --install-all || log "Falha ao recompilar módulos VMware (ignorado)"
    else
      log "VMware não instalado, ignorando recompilação de módulos."
    fi

    if [ -x /bin/desmonta_home.sh ]; then
      log "Desmontando home local com desmonta_home.sh..."
      sleep 10
      /bin/desmonta_home.sh
    else
      log "Script /bin/desmonta_home.sh não encontrado ou não executável."
    fi

    /bin/executa.sh
    ;;
  *)
    echo "Uso: /etc/init.d/ltsp-ssh-init start"
    exit 1
    ;;
esac